# gleipnird/gleipnird.service
[Unit]
Description=Gleipnir Daemon

[Service]
ExecStart=/usr/bin/gleipnird

[Install]
WantedBy=multi-user.target
